# voxbulk-api/supervisor/voxbulk.conf
; VoxBulk supervisor programs — install with scripts/vps-setup-celery.sh
; Replaces legacy retover-celery (bare "celery" not on PATH → FATAL).

[program:voxbulk-celery]
directory=/www/voxbulk/voxbulk-api
command=/www/voxbulk/scripts/run-celery-worker.sh
autostart=true
autorestart=true
startretries=5
startsecs=3
stopwaitsecs=30
stdout_logfile=/tmp/voxbulk-celery.log
stderr_logfile=/tmp/voxbulk-celery.err.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
environment=PYTHONUNBUFFERED="1"
